# web2py/myapp/cron/crontab
#crontab
@reboot  root *applications/myapp/cron/tablebar.py
